# config/crontab
* * * * * baca.py
#0 * * * * apiSend.py
#4,8,12 * * * * retryApiSend.py

