# config/crontab
* * * * * baca.py
#0 * * * * sendApi.py
